# curator4-1273683266.conf
[program:curator4-1273683266]
command=python3 /home/curator4/curator4-1273683266/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1273683266/test.err.log
stdout_logfile=/home/curator4/curator4-1273683266/test.out.log